# build/packages/packages/bsp/common/etc/update-motd.d/10-armbian-header
#!/bin/bash
#
# Copyright (c) Authors: http://www.armbian.com/authors
#
# This file is licensed under the terms of the GNU General Public
# License version 2. This program is licensed "as is" without any
# warranty of any kind, whether express or implied.

# DO NOT EDIT THIS FILE but add config options to /etc/default/armbian-motd
# any changes will be lost on board support package update

THIS_SCRIPT="header"
MOTD_DISABLE=""

[[ -f /etc/default/armbian-motd ]] && . /etc/default/armbian-motd

for f in $MOTD_DISABLE; do
	[[ $f == $THIS_SCRIPT ]] && exit 0
done

. /etc/os-release
. /etc/armbian-release

KERNELID=$(uname -r)

TERM=linux toilet -f standard -F metal "HID iClass SE"
TERM=linux toilet -f standard -F metal "RB25F"

#printf '\nWelcome to \e[0;91mRB25F\x1B[0m %s %s %s %s\n' "$VERSION $IMAGE_TYPE $PRETTY_NAME $KERNELID"
printf '\nWelcome to \e[0;91mHID iClass SE RB25F\x1B[0m %s %s %s 4.14.59\n' "$VERSION $IMAGE_TYPE $PRETTY_NAME"

